# templates/default/aptly-api.conf.erb
# aptly server upstart script

start on (net-device-up and local-filesystems and runlevel [2345])
stop on runlevel [016]

# Automatically restart process if crashed
respawn

script
  echo $$ > /var/run/aptly-api.pid
  exec aptly api serve -listen=:<%= node['aptly']['api']['port'] %>
end script

post-stop script
  rm -f /var/run/aptly-api.pid
end script
